# systemd/client.service.in
[Unit]
Description=User Mode HTTP Client

[Service]
WorkingDirectory=@WORKING_DIRECTORY@
ExecStart=@PYTHON3@ client.py --oneshot --server-url http://@SERVER_ADDR@:@PORT@/submit
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
