# systemd/server.service.in
[Unit]
Description=User Mode HTTP Server

[Service]
WorkingDirectory=@WORKING_DIRECTORY@
ExecStart=/usr/bin/python3 server.py -P @PORT@
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
